# service/spillman-paging.service
[Unit]
Description=Spillman Automation - Paging
After=syslog.target network.target

[Service]
WorkingDirectory=/opt/spillman-automation
ExecStart=/usr/bin/python3 /opt/spillman-automation/app.py --paging

Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target